# netscout.timer
[Unit]
Description="tasteful software downloader timer"
[Timer]
OnBootSec=15min
OnActiveSec==6hr

[Install]
WantedBy=multi-user.target
