# netscout.timer
[Unit]
Description="tasteful software downloader timer"
[Timer]
OnCalendar=6hr
[Install]
WantedBy=multi-user.target
